# dehydrated/files/etc/cron.weekly/dehydrated
#!/bin/bash
sudo -u dehydrated /home/dehydrated/dehydrated/dehydrated --cron
